# scripts/usb-nic-monitor.service
[Unit]
Description=USB NIC Stability Monitor for k3s1 and k3s2
After=network.target

[Service]
Type=simple
ExecStart=/home/josh/flux-k3s/scripts/monitor-usb-nic.sh
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=usb-nic-monitor

[Install]
WantedBy=multi-user.target